# docker/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=rdfmap
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
priority=10

[program:api]
command=/usr/local/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.err
priority=20
environment=PYTHONPATH="/app/backend"

[program:worker]
command=/usr/local/bin/celery -A app.worker:celery_app worker --loglevel=info --concurrency=2
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker.err
priority=30
environment=PYTHONPATH="/app/backend"

